# docker-context/cron.d/log-reboot
@reboot root echo "in crond after reboot"
