# docker-context/cron.d/timezone
@reboot root /bin/timedatectl set-timezone Asia/Kolkata