# docker-context/cron.d/tmpfs-ramdisk
@reboot root /usr/local/bin/mount-ramdisk.sh
